# container/chatterbox-tts.container
# Chatterbox TTS Quadlet Container Unit (CPU)
#
# This file defines the Chatterbox TTS container for systemd/Podman Quadlet.
# Uses socket activation - the container starts when a client connects.
#
# The client auto-detects GPU and uses the CUDA backend if available.
# Install both CPU and CUDA versions for automatic fallback.
#
# Installation (user):
#   mkdir -p ~/.config/containers/systemd
#   cp chatterbox-tts.socket ~/.config/containers/systemd/
#   cp chatterbox-tts.container ~/.config/containers/systemd/
#   systemctl --user daemon-reload
#   systemctl --user enable --now chatterbox-tts.socket
#
# Installation (system):
#   cp chatterbox-tts.socket /etc/containers/systemd/
#   cp chatterbox-tts.container /etc/containers/systemd/
#   systemctl daemon-reload
#   systemctl enable --now chatterbox-tts.socket
#

[Unit]
Description=Chatterbox TTS Daemon for Speech Dispatcher
Requires=chatterbox-tts.socket
After=chatterbox-tts.socket

[Container]
# Container image from GitHub Container Registry
Image=ghcr.io/rsturla/chatterbox-spd:latest
# For local development builds, use:
# Image=localhost/chatterbox-tts:latest

# Container name
ContainerName=chatterbox-tts

# Run parameters for the daemon
# Socket is inherited via FD from systemd (socket activation)
Exec=--voice-dir /app/voices --output-dir /app/output

# No network needed - uses socket activation
Network=none

# Audio output directory - shared with host for playback
Volume=%t/chatterbox-tts:/app/output:z

# Voice files - for voice cloning
Volume=%h/.cache/chatterbox-spd/voices:/app/voices:z

# Models are baked into the image - no cache volume needed

# Environment
Environment=PYTHONUNBUFFERED=1
# Offline mode - models baked into image, dummy token satisfies token=True check
Environment=HF_HUB_OFFLINE=1
Environment=HF_TOKEN=offline_mode_dummy_token
# Idle timeout in seconds (0 to disable, default: 1200 = 20 minutes)
Environment=CHATTERBOX_IDLE_TIMEOUT=1200

# Auto-update policy
AutoUpdate=registry

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=900

[Install]
# Socket activates this service, not default.target
WantedBy=
