# container/chatterbox-tts-cuda.container
# Chatterbox TTS Quadlet Container Unit (CUDA/GPU)
#
# This file defines the Chatterbox TTS container with GPU support.
# Use this version if you have an NVIDIA GPU.
# Uses socket activation - the container starts when a client connects.
#
# Requirements:
#   - NVIDIA GPU
#   - nvidia-container-toolkit installed
#   - Podman configured for GPU access
#
# Installation (user):
#   mkdir -p ~/.config/containers/systemd
#   cp chatterbox-tts-cuda.socket ~/.config/containers/systemd/
#   cp chatterbox-tts-cuda.container ~/.config/containers/systemd/
#   systemctl --user daemon-reload
#   systemctl --user enable --now chatterbox-tts-cuda.socket
#
# The client will auto-detect GPU and use this service when available.
#

[Unit]
Description=Chatterbox TTS Daemon for Speech Dispatcher (CUDA)
Requires=chatterbox-tts-cuda.socket
After=chatterbox-tts-cuda.socket

[Container]
# Container image from GitHub Container Registry (CUDA version)
Image=ghcr.io/rsturla/chatterbox-spd:cuda
# For local development builds, use:
# Image=localhost/chatterbox-tts:cuda

# Container name
ContainerName=chatterbox-tts-cuda

# Run parameters for the daemon
# Socket is inherited via FD from systemd (socket activation)
Exec=--voice-dir /app/voices --output-dir /app/output --device cuda

# No network needed - uses socket activation
Network=none

# GPU access - pass all NVIDIA devices
AddDevice=nvidia.com/gpu=all

# Security - required for GPU access
SecurityLabelDisable=true

# Audio output directory - shared with host for playback
Volume=%t/chatterbox-tts-cuda:/app/output:z

# Voice files - for voice cloning
Volume=%h/.cache/chatterbox-spd/voices:/app/voices:z

# Models are baked into the image - no cache volume needed

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=NVIDIA_VISIBLE_DEVICES=all
Environment=NVIDIA_DRIVER_CAPABILITIES=compute,utility
# Offline mode - models baked into image, dummy token satisfies token=True check
Environment=HF_HUB_OFFLINE=1
Environment=HF_TOKEN=offline_mode_dummy_token
# Idle timeout in seconds (0 to disable, default: 1200 = 20 minutes)
Environment=CHATTERBOX_IDLE_TIMEOUT=1200

# Auto-update policy
AutoUpdate=registry

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=300

[Install]
# Socket activates this service, not default.target
WantedBy=
